# packages/wayland/compositor/sway/system.d/sway-monitor-toggle.service
[Unit]
Description=Handle monitor plug events

[Service]
Type=simple
ExecStart=/usr/bin/sway-monitor-plugged.sh

[Install]
WantedBy=multi-user.target
