# packages/wayland/compositor/sway/system.d/sway-monitor-toggle.service
[Unit]
Description=Handle monitor plug events
After=essway.service

[Service]
Type=simple
User=root
Environment=HOME=/storage
Environment=XDG_RUNTIME_DIR=/var/run/0-runtime-dir
ExecStart=/usr/bin/sway-monitor-plugged.sh

[Install]
WantedBy=rocknix.target
