# packages/wayland/compositor/sway/udev.d/99-on-monitor-plugged.rules
KERNEL=="card*", SUBSYSTEM=="drm", ACTION=="add|change", RUN+="/usr/bin/systemctl start sway-monitor-toggle.service"
